# scripts/thinkpad-power/99-thinkpad-mute-led.rules
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="platform::mute", ATTR{trigger}="none", RUN+="/bin/chmod 0666 /sys/class/leds/platform::mute/brightness"
