# scripts/thinkpad-power/99-thinkpad-profile.rules
# Prevent firmware from resetting platform profile on AC plug/unplug
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/bin/thinkpad-profile-hold"
